# PowerSwitcher/60-onbattery.rules
SUBSYSTEM=="power_supply",ENV{POWER_SUPPLY_ONLINE}=="0",RUN+="/usr/local/bin/powerswitcher.sh bal"
